# examples/cron/sensos-server-backup.cron.example
MAILTO=""
PATH=/usr/local/bin:/usr/bin:/bin

# Run server backups every day at 02:15.
15 2 * * * cd /home/sensos/sensos-server && ./bin/backup-server.sh >> /home/sensos/sensos-server/local/log/backup-cron.log 2>&1

# Example with built-in rclone export to Box:
# 15 2 * * * cd /home/sensos/sensos-server && ./bin/backup-server.sh --export --remote box:sensos-server-backups >> /home/sensos/sensos-server/local/log/backup-cron.log 2>&1

# Example with a local post-hook:
# 15 2 * * * cd /home/sensos/sensos-server && ./bin/backup-server.sh --post-hook ./local/hooks/post-backup.sh >> /home/sensos/sensos-server/local/log/backup-cron.log 2>&1
